# cmake/catcierge.service.in
[Unit]
Description=Catcierge automatic cat door awesomeness
@CATCIERGE_SERVICE_AFTER@

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/catcierge_grabber
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console
User=catcierge
Group=users

[Install]
WantedBy=multi-user.target
